# roles/muffins_cache_mover/templates/cache_mover.timer.j2
{{ ansible_managed | comment }}
[Unit]
Description=Timer for cache_mover service

[Timer]
OnCalendar={{ on_calendar_schedule }}
Persistent=true

[Install]
WantedBy=timers.target
